# systemd/user/poll-email@.timer
[Unit]
Description=Poll email timer for account %i

[Timer]
OnBootSec=1m
OnUnitInactiveSec=15

[Install]
WantedBy=default.target
