# systemd/user/sync-email@.timer
[Unit]
Description=Sync email timer for account %i

[Timer]
OnBootSec=1m
OnUnitInactiveSec=5m

[Install]
WantedBy=default.target
